# crontab
# Raspored prikupljanja (vrijeme = lokalno, Europe/Zagreb zbog TZ u Dockerfileu).
# Cilj: svakih 30 min, od 5:30 do 23:00, SVAKI dan. Oba smjera (2 rute).
# Tri retka da granice budu točno 5:30 i 23:00 (a ne 5:00 / 23:30).
#
# Procjena: ~36 termina/dan x 2 smjera x 31 = ~2230 poziva/mjesec (< 5000 free).
30 5    * * * cd /app && python collect.py
0,30 6-22 * * * cd /app && python collect.py
0 23    * * * cd /app && python collect.py
